# app.ini
[program:dashboard]
command = /usr/bin/python3 %(APP_PATH)s/dashboard.py
autostart = true
autorestart = true
startsecs = 10
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 100
stdout_capture_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 100
stderr_capture_maxbytes = 10MB
stdout_logfile = %(APP_PATH)s/tmp/stdout.log
stderr_logfile = %(APP_PATH)s/tmp/stderr.log
